# server/docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:php-fpm]
command=/etc/init.d/php8.2-fpm start -F
autostart=true
autorestart=false
startretries=1
numprocs=1
startsecs=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=false
startretries=1
numprocs=1
startsecs=0

[program:composer]
command=composer install
autostart=true
autorestart=false
startretries=1
numprocs=1
startsecs=0
